# install/supervisord.conf
[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0777

[supervisord]
pidfile=/var/run/supervisord.pid
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
directory = /tmp
nocleanup = true
childlogdir = /tmp

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock


[program:domopyc]
command=/home/pi/venv/bin/python /home/pi/venv/lib/python3.4/site-packages/domopyc/domopyc_main.py
process_name=%(program_name)s
numprocs=1

